# src/services.d/odoo/finish
#!/bin/sh

if [ ! -z "$ONESHOT" ]; then
    exec s6-svscanctl -t /var/run/s6/services
fi