# src/services.d/odoo/finish
#!/bin/sh

if [ ! -z "$ODOO_ONESHOT" ]; then
    /run/s6/basedir/bin/halt
fi