# src/services.d/odootail/run
#!/usr/bin/env bash

cd /opt/odoo

tail -f "$(odooconf options logfile)"