# src/services.d/odootail/run
#!/usr/bin/env bash

cd /opt/odoo

# Check if ODOOCONF__options__logfile is set
if [ ! -z "$ODOOCONF__options__logfile" ]; then
    tail -f "$ODOOCONF__options__logfile"
else
    sleep infinity
fi